# server/systemd/daily-shutdown.service
[Unit]
Description=2 AM to 1030 AM

[Service]
ExecStart=/usr/sbin/shutdown -h now
