# server/systemd/user/doctor-check.service
[Unit]
Description=Daily server health check

[Service]
Type=oneshot
ExecStart=%h/bin/doctor-check
